# deployment/prayer-player.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prayerplayer.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/PrayerPlayer.app/Contents/MacOS/PrayerPlayer</string>
        <!-- Add any arguments your script needs here, e.g., --city "London" -->
        <!-- <string>--city</string> -->
        <!-- <string>YourCity</string> -->
        <!-- <string>--country</string> -->
        <!-- <string>YourCountry</string> -->
    </array>

    <!-- Run at load (when the user logs in or the agent is loaded) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- KeepAlive: If true, launchd will restart the job if it exits. -->
    <key>KeepAlive</key>
    <true/>

    <!-- StandardOutPath: Path to redirect standard output -->
    <key>StandardOutPath</key>
    <string>/tmp/com.prayerplayer.scheduler.stdout.log</string>

    <!-- StandardErrorPath: Path to redirect standard error -->
    <key>StandardErrorPath</key>
    <string>/tmp/com.prayerplayer.scheduler.stderr.log</string>

    

    <!-- EnvironmentVariables: Set environment variables for the script -->
    <!-- <key>EnvironmentVariables> -->
    <!-- <dict> -->
    <!--     <key>PYTHONPATH</key> -->
    <!--     <string>/path/to/your/prayer/src</string> -->
    <!-- </dict> -->

</dict>
</plist>